# rm
!/bin/bash
### BEGIN INIT INFO
# Provides: MyService
# Required-Start:    $all
# Required-Stop: 
# Default-Start:     5 
# Default-Stop:      6 
# Short-Description: Monitoring by Pavlovsk DOC 
### END INIT INFO

/home/pi/Documents/PythonScripts/raspi_doc/main.py
